# profile.d/mysql.sh
# Update PATH to include @@PREFIX@@/bin

PATH="@@PREFIX@@/bin:$PATH"

export PATH
