# profile.d/postfix.sh
# Update PATH to include @@PREFIX@@/sbin and @@PREFIX@@/bin

PATH="@@PREFIX@@/bin:@@PREFIX@@/sbin:$PATH"

export PATH
